# priv/automation/macos/dev.cataladev.leetcode-sync.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>dev.cataladev.leetcode-sync</string>

  <key>ProgramArguments</key>
  <array>
    <string>/Users/carlos/Development/Leetcode Auto Answer Uploader/scripts/run_sync.sh</string>
  </array>

  <key>WorkingDirectory</key>
  <string>/Users/carlos/Development/Leetcode Auto Answer Uploader</string>

  <key>RunAtLoad</key>
  <true/>

  <key>StartCalendarInterval</key>
  <dict>
    <key>Hour</key>
    <integer>9</integer>
    <key>Minute</key>
    <integer>0</integer>
  </dict>

  <key>StandardOutPath</key>
  <string>/Users/carlos/Development/Leetcode Auto Answer Uploader/logs/launchd.out.log</string>

  <key>StandardErrorPath</key>
  <string>/Users/carlos/Development/Leetcode Auto Answer Uploader/logs/launchd.err.log</string>
</dict>
</plist>
